# dashboard/senpi-dashboard.service
[Unit]
Description=Senpi Trading Dashboard
After=network.target

[Service]
Type=simple
User=senpi
Group=senpi
WorkingDirectory=/opt/senpi/senpi-waifu/dashboard
Environment=SENPI_WAIFU_DIR=/opt/senpi/senpi-waifu
Environment=DASH_PORT=8420
EnvironmentFile=/opt/senpi/.env
ExecStart=/opt/senpi/venv/bin/uvicorn server:app --host 0.0.0.0 --port 8420
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
